# tagesschau-schauer.timer
[Unit]
Description=Sends notification from Tagesschau news feed

[Timer]
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
